# UniCoreLinux/build/archiso/airootfs/etc/systemd/system/kextd.service
[Unit]
Description=UniCore kext loader (kextd)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kextd auto
Restart=on-failure

[Install]
WantedBy=multi-user.target
