# UniCoreLinux/build/archiso/airootfs/etc/systemd/system/session-restore.service
[Unit]
Description=UniCore Session Restore
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/session-restore restore
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
